# projects/services/jenkins/conf/supervisor/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock   
chmod=0700                       

[supervisord]
logfile=/var/log/supervisor/supervisord.log 
pidfile=/var/run/supervisord.pid 
childlogdir=/var/log/supervisor            

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock 

[supervisord]
loglevel=debug
nodaemon=true

[program:jenkins]
environment = JENKINS_HOME="/var/jenkins_home", JENKINS_PREFIX="/"
command = /usr/local/bin/jenkins.sh

[program:docker]
command = docker -d
autostart = true
autorestart = false

[include]
files = /etc/supervisor/conf.d/*.conf
